# ArcherC7v2/etc/crontabs/root
0 1,3,5 * * * wifi down
1 7 * * * reboot
